# config/systemd/user/hyprsunset-stop.timer
[Unit]
Description=hyprsunset at 11pm

[Timer]
OnCalendar=Mon,Tue,Wed,Thu,Fri 8:00

[Install]
WantedBy=timers.target
